# scripts/voice/com.deepstack.voice-listener.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.deepstack.voice-listener</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/eddiebelaval/Development/id8/products/deepstack/scripts/voice/deepstack-voice-listener.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/Users/eddiebelaval/Development/id8/products/deepstack/venv/bin:/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
        <key>HOME</key>
        <string>/Users/eddiebelaval</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/eddiebelaval/Development/id8/products/deepstack</string>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>10</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>Nice</key>
    <integer>15</integer>

    <key>LowPriorityIO</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/eddiebelaval/Development/id8/products/deepstack/logs/voice-listener-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/eddiebelaval/Development/id8/products/deepstack/logs/voice-listener-stderr.log</string>
</dict>
</plist>
